# twemcache/DEBIAN/postinst
#!/bin/sh
update-rc.d -f twemcache defaults
