# twemcache/DEBIAN/prerm
#!/bin/sh
update-rc.d -f twemcache remove
